# homeamp-v2/deployment/homeamp-api.service
[Unit]
Description=HomeAMP V2.0 API Service
After=network.target mariadb.service
Requires=mariadb.service

[Service]
Type=simple
User=webadmin
Group=webadmin
WorkingDirectory=/opt/homeamp/homeamp.ampdata/homeamp-v2
ExecStart=/opt/homeamp/venv/bin/python -m uvicorn homeamp_v2.api.main:app --host 0.0.0.0 --port 8077 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homeamp-api

[Install]
WantedBy=multi-user.target
