# homeamp-v2/deployment/homeamp-agent.service
[Unit]
Description=HomeAMP V2.0 Discovery Agent
After=network.target mariadb.service homeamp-api.service
Requires=mariadb.service

[Service]
Type=simple
User=webadmin
Group=webadmin
WorkingDirectory=/opt/homeamp/homeamp.ampdata/homeamp-v2
ExecStart=/usr/bin/python3 -m homeamp_v2.agent
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homeamp-agent

[Install]
WantedBy=multi-user.target
