# software/homeamp-config-manager/deployment/plugin-update-scheduler/archivesmp-plugin-update.timer
[Unit]
Description=Hourly ArchiveSMP Plugin Update Check
Requires=archivesmp-plugin-update.service

[Timer]
OnCalendar=hourly
Persistent=true
Unit=archivesmp-plugin-update.service

[Install]
WantedBy=timers.target
